# base/roles/armbian_base/templates/btrbk.timer.j2
[Unit]
Description=btrbk periodic snapshot

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=5min

[Install]
WantedBy=timers.target
